# etc/cron/cron.hourly/seedboxsync
#!/bin/sh

/usr/bin/seedboxsync-launcher
